# systemd/kaosgdd-frontend.service
[Unit]
Description=KaosGdd Web Frontend
After=network.target kaosgdd-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/srv/KaosGdd-web/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=NEXT_PUBLIC_API_BASE=http://127.0.0.1:8000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
