# systemd/kaosgdd-lifecycle.service
[Unit]
Description=KaosGdd Lifecycle Maintenance Sweep
After=kaosgdd-backend.service
Requires=kaosgdd-backend.service

[Service]
Type=oneshot
ExecStart=/usr/bin/curl --fail --show-error --silent -X POST http://127.0.0.1:8000/internal/lifecycle/maintain
