# systemd/kaosgdd-reminders-scan-missed.service
[Unit]
Description=KaosGdd Scan Missed Reminders
After=kaosgdd-backend.service
Requires=kaosgdd-backend.service

[Service]
Type=oneshot
ExecStart=/usr/bin/curl --fail --show-error --silent -X POST http://127.0.0.1:8000/internal/reminders/scan-missed
